# systemd/irexec.service
[Unit]
Documentation=man:irexec(1)
Documentation=http://lirc.org/html/configure.html
Documentation=http://lirc.org/html/configure.html#lircrc_format
Description=Handle events from IR remotes decoded by lircd(8)

[Service]
; user =
; group =
Type=simple
ExecStart=/usr/sbin/irexec /etc/lirc/irexec.lircrc

[Install]
WantedBy=multi-user.target
